# deploy/unation.service
# ============================================================
# unation.service — systemd 서비스 파일 (Linux VPS용)
#
# 이 파일을 쓰면 서버가 재부팅되거나 죽어도 자동으로 다시 시작되어
# 24시간 365일 운영이 가능합니다. (가비아/Cafe24 VPS, AWS Lightsail 등
# 우분투/데비안 계열 리눅스 서버 기준)
#
# 설치 방법 (deployment-guide.md 3단계 참고):
#   1. 아래 User, WorkingDirectory 경로를 본인 서버에 맞게 수정
#   2. sudo cp unation.service /etc/systemd/system/unation.service
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable unation      (재부팅 시 자동 시작)
#   5. sudo systemctl start unation       (지금 바로 시작)
#
# 상태 확인: sudo systemctl status unation
# 로그 확인: sudo journalctl -u unation -f
# 재시작:    sudo systemctl restart unation
# ============================================================

[Unit]
Description=UNATION (유네이션) Node.js Server
After=network.target

[Service]
Type=simple

# 실제 서버에 맞게 수정하세요 (예: ubuntu, 또는 별도로 만든 unation 전용 사용자)
User=ubuntu

# 프로젝트를 업로드한 실제 경로로 수정하세요
WorkingDirectory=/home/ubuntu/unation

# node 실행 파일의 전체 경로를 모를 경우 `which node` 명령으로 확인 후 수정하세요
ExecStart=/usr/bin/node /home/ubuntu/unation/server/index.js

# 서버가 죽으면 5초 후 자동 재시작 (정전/일시적 오류 등에서도 자동 복구)
Restart=always
RestartSec=5

# .env 파일이 같은 폴더에 있으면 server/env.js가 자동으로 읽어들입니다.
# systemd 레벨에서 환경변수를 직접 지정하고 싶다면 아래 주석을 풀고 사용하세요.
# Environment=NODE_ENV=production
# Environment=PORT=3000
# EnvironmentFile=/home/ubuntu/unation/.env

# 표준출력/에러를 시스템 저널로 보냄 (journalctl -u unation 으로 확인)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
